# build/healarr.conf
[program:healarr-script]
user = nobody
startsecs = 0
autorestart = False
command = /home/nobody/healarr.sh
umask = 000
